# bin/dht11_nodejs_socket.service
[Unit]
Description=DHT11 on express with socket.io

[Service]
Type=simple
ExecStart=/home/pi/dht11_nodejs_socket/bin/www
#Restart=on-failure
Restart=always
RestartSec=10  
#User=pi
#Group=pi
Environment=PATH=/usr/bin:/usr/local/bin
Environment=NODE_ENV=production
WorkingDirectory=/home/pi/dht11_nodejs_socket

[Install]
WantedBy=multi-user.target
